# scripts/crontab/install.sh
#!/bin/sh 
cp -f /vagrant/scripts/crontab/laravel-task /etc/cron.d/laravel-task