# scripts/crontab/laravel-task
 * *   *   *   *   vagrant /usr/bin/php /home/vagrant/Code/LA1-Website/artisan scheduled:run