# scripts/crontab/laravel-task
 * *   *   *   *   vagrant cd /home/vagrant/Code/LA1-Website && /usr/bin/php artisan scheduled:run